# infra/vpsautodeploy.timer
[Unit]
Description=Roda o auto-deploy a cada 30 segundos

[Timer]
OnBootSec=90
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
