# infra/vpsbackup.timer
[Unit]
Description=VPS Backup (relogio minuto a minuto - perfis decidem)

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target
